# com.lewis.jobhunt.telegram-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lewis.jobhunt.telegram-poller</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/lewis/Desktop/agent/venv/bin/python3</string>
    <string>/Users/lewis/Desktop/agent/src/api/telegram_poller.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/lewis/Desktop/agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/lewis/Desktop/agent/src</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/lewis/Desktop/agent/logs/telegram_poller.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lewis/Desktop/agent/logs/telegram_poller.launchd.err.log</string>
</dict>
</plist>
